# target/linux/bcm53xx/base-files/etc/uci-defaults/00-enable-services
/etc/init.d/rpcd enable
/etc/init.d/ubus-scriptd enable
/etc/init.d/uhttpd enable
